# controller/systemd/playground-controller.service
# Step 3.7 systemd unit for the playground session controller.
#
# UNTESTED-PENDING-HOST: this file has not been installed or started on a
# real CT 105 — it was written against the plan's requirements (restart on
# failure, token loaded from an environment file, not hardcoded) and has not
# been verified with `systemctl status` / a kill-and-restart test, both of
# which are explicit Step 3.7 verification gates.
#
# Deploy layout assumed here (adjust paths if Step 3.7 deploys differently):
#   /opt/playground-controller/           <- this repo's controller/ dir
#   /opt/playground-controller/venv/      <- Python 3.11+ virtualenv
#   /etc/playground-controller/controller.env  <- copied from .env.example,
#       filled in with the real token from Step 3.3, chmod 600

[Unit]
Description=Playground session controller (Phase 3)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=playground-ctrl
Group=playground-ctrl
WorkingDirectory=/opt/playground-controller
EnvironmentFile=/etc/playground-controller/controller.env
ExecStart=/opt/playground-controller/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
